# ksnapshotservice.service.in
[Unit]
Description=Snapshot access service for KDE

[Service]
Type=dbus
BusName=org.kde.ksnapshotservice
ExecStart=@KDE_INSTALL_FULL_LIBEXECDIR@/ksnapshotservice
